# packaging/systemd/wol-libvirt-gateway.service
[Unit]
Description=Wake-on-LAN to libvirt gateway
Documentation=https://github.com/brunoproduit/wol-libvirt-gateway
After=network.target libvirtd.service
Wants=libvirtd.service
BindsTo=libvirtd.service

[Service]
Type=simple
ExecStart=/usr/bin/wol-libvirt-gateway --address 0.0.0.0:9 --libvirt-uri qemu:///system
Restart=always
RestartSec=5
User=wol-libvirt-gateway
Group=wol-libvirt-gateway

# Allow access to libvirt socket
ReadWritePaths=/var/run/libvirt/libvirt-sock
ReadWritePaths=/var/run/libvirt/libvirt-sock-ro

# Need NET_BIND_SERVICE for binding to port 9 (< 1024)
CapabilityBoundingSet=CAP_NET_BIND_SERVICE
AmbientCapabilities=CAP_NET_BIND_SERVICE

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProtectProc=invisible
ProcSubset=pid
PrivateTmp=true
PrivateDevices=true
RemoveIPC=true
LockPersonality=true
RestrictRealtime=true
RestrictSUIDSGID=true
RestrictNamespaces=true
MemoryDenyWriteExecute=true
SystemCallArchitectures=native
SystemCallFilter=@system-service

# Network restrictions
IPAddressDeny=any
IPAddressAllow=127.0.0.1/8

# Resource limits
LimitNOFILE=1024
LimitNPROC=64

# Environment
Environment=RUST_LOG=info

[Install]
WantedBy=multi-user.target
